# systemd/cosmic-nightlight.service
# SPDX-License-Identifier: MPL-2.0
#
# User-level systemd service for cosmic-nightlight's background daemon.
#
# Install with:
#   mkdir -p ~/.config/systemd/user
#   cp systemd/cosmic-nightlight.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now cosmic-nightlight.service
#
# Adjust ExecStart below if the binary isn't installed to /usr/local/bin.

[Unit]
Description=Cosmic Night Light background service
After=graphical-session.target

[Service]
ExecStart=/usr/local/bin/cosmic-nightlight --daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
