# debian/cosmic-nightshift.user.service
# SPDX-License-Identifier: MPL-2.0
# Installed to /usr/lib/systemd/user/cosmic-nightshift.service by
# dh_installsystemduser. Enable per-user with:
#   systemctl --user enable --now cosmic-nightshift.service

[Unit]
Description=COSMIC Night Shift background scheduler
After=graphical-session.target
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/cosmic-nightshift --daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
